# GridboxConnectorAddon-experimental/rootfs/etc/s6-overlay/s6-rc.d/gridboxconnector/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Viessmann Gridbox Connector
# ==============================================================================

bashio::log.info "Starting GridboxConnector..."
bashio::log.debug "lsusr: ${ls -lash /usr/bin/}"
bashio::log.debug "lsbuild: ${ls -lash /build/}"
exec /usr/bin/gridboxconnector.sh